# debian/cron.d
*/5 * * * *     root    /usr/bin/redmine-mail2ticket

                                                                                                         