# debian/redmine-mail2ticket.cron.d
*/5 * * * *     root    /usr/share/redmine/extra/mail_handler/mail2ticket.rb
                                                                                                      